# usr/etc/profile.d/aurora-fastfetch.sh
alias fastfetch="/usr/bin/fastfetch --file /usr/share/ublue-os/aurora-logo.txt